# .config/systemd/system/update-system.timer
[Unit]
Description=Update system weekly

[Timer]
OnCalendar=weekly
Persistent=true 

[Install]
WantedBy=timers.target
